# config/supervisor/influxdb.conf
[program:influxdb]
command=/opt/influxdb/influxd -config /etc/opt/influxdb/influxdb.conf
autostart=true
autorestart=true
user=root
priority=1
