# config/supervisor/rtcscope-server.conf
[program:rtcscopeserver]
command=/usr/bin/python app.py
directory=/opt/rtcscope-server
autostart=true
autorestart=true
user=root
priority=3
